# ops/systemd/eventicro-api.service
[Unit]
Description=EventiCRO External API
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=eventicro
Group=eventicro
WorkingDirectory=/opt/eventicro-mono/apps/api
Environment=NODE_ENV=production
Environment=PORT=8787
EnvironmentFile=/etc/default/eventicro.env
ExecStart=/usr/bin/node dist/index.js
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=eventicro-api

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/eventicro-mono/apps/api/logs

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
