# ops/systemd/eventicro-api-preview.service
[Unit]
Description=EventiCRO API preview (Node on 127.0.0.1:8787)
After=network.target

[Service]
Type=simple
User=eventicro-api
WorkingDirectory=/srv/app/src
Environment=HOST=127.0.0.1
Environment=PORT=8787
Environment=NODE_ENV=production
ExecStart=/usr/bin/env node /srv/app/src/scripts/preview.js
Restart=always
RestartSec=2
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
