# ops/systemd/eventicro-ui.service
[Unit]
Description=EventiCRO UI (SvelteKit SSR on 127.0.0.1:3000)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/srv/app/ui
ExecStart=/usr/local/bin/run-eventicro-ui
Restart=always
RestartSec=3
# CRITICAL: Force the app to listen on the correct IP
Environment="HOST=127.0.0.1"
Environment="PORT=3000"
Environment="HOME=/root"

[Install]
WantedBy=multi-user.target
