# ops/systemd/eventicro-ui.service
[Unit]
Description=EventiCRO UI (SvelteKit SSR on 127.0.0.1:3000)
After=network-online.target
Wants=network-online.target

[Service]
Environment=HOME=/root
Type=simple
# User and HOME are critical for nvm to work correctly
User=root
Environment=HOME=/root

# Load DB credentials etc.
EnvironmentFile=-/etc/default/eventicro.env

# Runtime defaults (wrapper can override if needed)
Environment=HOST=127.0.0.1
Environment=PORT=3000
Environment=NODE_ENV=production

WorkingDirectory=/srv/app/ui
ExecStart=/usr/local/bin/run-eventicro-ui
Restart=always
RestartSec=2
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
